# bot-multi-gold.service
[Unit]
Description=Bot Multi Gold - 8-Indicator Trading Bot (XAUUSD)
After=network.target
Wants=network.target

[Service]
Type=simple
User=root
WorkingDirectory=/opt/bot_multi_gold
ExecStart=/usr/bin/python3 /opt/bot_multi_gold/main.py
Restart=always
RestartSec=10
StandardOutput=journal
StandardError=journal

# Environment variables
Environment=PYTHONPATH=/opt/bot_multi_gold
Environment=PYTHONUNBUFFERED=1

# Security settings (basitleştirilmiş)
NoNewPrivileges=true

[Install]
WantedBy=multi-user.target
